# package/mtk/applications/hnat-detect/files/etc/hotplug.d/iface/99-hnat-detect
#!/bin/sh
[ "$ACTION" = "ifup" ] || [ "$ACTION" = "update" ] || exit 0
[ "$INTERFACE" = "loopback" ] && exit 0

/usr/bin/ucode /usr/share/ucode/hnat/detect.uc
